# sws-archive.timer
[Unit]
Description=Archive meteodata.log periodically

[Timer]
OnCalendar=Monthly

[Install]
WantedBy=timers.target
